# backend/deploy/gr-api.service
[Unit]
Description=GR API (FastAPI + Uvicorn)
After=network.target postgresql.service

[Service]
Type=simple
User=app_gr
Group=app_gr
WorkingDirectory=/opt/app_gr/backend
EnvironmentFile=/opt/app_gr/backend/.env
# --workers 1: cada worker uvicorn instancia seu proprio APScheduler
# (lifespan roda por processo). Com workers=2 cada endpoint disparava 2x
# por tick, dobrando chamadas a fontes externas (ex.: QiTech). Decisao
# 2026-05-09: 1 worker e suficiente pra carga atual; quando precisar de
# paralelismo HTTP, separar processo do scheduler do processo HTTP em
# vez de aumentar workers aqui. Detalhe: com `--workers 1`, uvicorn nao
# faz fork — main process eh o worker tambem (sem multiprocessing).
# --host 0.0.0.0: o TLS de callback.strataai.com.br e terminado por um proxy
# EXTERNO a esta VM (nao ha :443 local) que alcanca o gr-api direto no IP de
# LAN (192.168.100.26:8000). Bind em 127.0.0.1 derruba TODO o trafego externo
# (login/webhooks) -> manter 0.0.0.0. O firewall/rede restringe o acesso.
ExecStart=/opt/app_gr/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 1
Restart=on-failure
RestartSec=3
StandardOutput=append:/var/log/gr-api/stdout.log
StandardError=append:/var/log/gr-api/stderr.log

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/var/log/gr-api /opt/app_gr/backend
ProtectHome=yes
# Storage de uploads de dossie (DOSSIER_STORAGE_ROOT=/var/lib/gr/dossier-attachments).
# StateDirectory cria /var/lib/gr (dono app_gr, 0755) E o adiciona ao conjunto
# gravavel do sandbox ProtectSystem=strict. Sem isto, gravar anexos/documentos
# de dossie estoura `OSError: Read-only file system` -> 500 sem header CORS no
# upload (o ServerErrorMiddleware fica fora do CORSMiddleware, entao o browser
# mascara como erro de CORS). Robusto a rebuild: garante a pasta em maquina nova.
StateDirectory=gr

[Install]
WantedBy=multi-user.target
